# skel/containers/nginx/app/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/app/start_nginx

[program:watch_and_reload]
command=/app/watch_and_reload
